# support/crontab
*/25 * * * *    cd /home/ubuntu/instabot && sh runner.sh >> app.log 2>&1
