# crontab
0 */6 * * *    cd /home/ubuntu/instabot && sh maintenance.sh >> maintenance.log 2>&1
